# triage/systemd/fleet-triage-filter.timer
[Unit]
Description=Run fleet triage filter every 5 minutes

[Timer]
OnBootSec=60s
OnUnitInactiveSec=5min
Unit=fleet-triage-filter.service
Persistent=true

[Install]
WantedBy=timers.target
